# service/xbit-match-adausdt.service
[Unit]
Description=MatchProcessDaemon ADAUSDT
After=network.target

[Service]
ExecStart=/usr/bin/python2 /data/bitstorm/daemon/MatchProcess/MatchProcessDaemon.py --product=ADAUSDT

User=centos
Group=centos

[Install]
WantedBy=multi-user.target
