# service/xbit-ordersend-bchusdt.service
[Unit]
Description=xbitOrdersend BCHUSDT
After=network.target

[Service]
ExecStart=/usr/bin/python2 /data/bitstorm/daemon/Operation/xbitOrdersend.py --product=BCHUSDT

User=centos
Group=centos

[Install]
WantedBy=multi-user.target
