# service/xbit-match-xrpusdt.service
[Unit]
Description=MatchProcessDaemon XRPUSDT
After=network.target

[Service]
ExecStart=/usr/bin/python2 /data/bitstorm/daemon/MatchProcess/MatchProcessDaemon.py --product=XRPUSDT

User=centos
Group=centos

[Install]
WantedBy=multi-user.target
