# service/xbit-match-xvgusdt.service
[Unit]
Description=MatchProcessDaemon XVGUSDT
After=network.target

[Service]
ExecStart=/usr/bin/python2 /data/bitstorm/daemon/MatchProcess/MatchProcessDaemon.py --product=XVGUSDT

User=centos
Group=centos

[Install]
WantedBy=multi-user.target